# defmath/files/supervisor/rastreio.conf
[program:rastreio]
command = /home/sam/.local/bin/gunicorn rastreio:app -b localhost:8000
directory = /home/sam/www/rastreio
user = sam